# Files/amd64/FHS/etc/local.d/01-bvfsom-initializer.start
#!/bin/bash
sudo "/System/Bit-OS Virtual File System Overlay Manager/bvfsom.sh"