# Files/amd64/FHS/etc/local.d/01-bvfsom-initializer.start
#!/bin/sh
if [ -f /usr/libexec/rc/sh/functions.sh ]; then
    . /usr/libexec/rc/sh/functions.sh
fi
sudo "/System/Bit-OS Virtual File System Overlay Manager/bvfsom.sh"